# docker/rootfs/etc/s6-overlay/s6-rc.d/cherry-log-prepare/up
if { mkdir -p /var/log/cherry }
if { chown nobody:nogroup /var/log/cherry }
chmod 02755 /var/log/cherry
