# docker/rootfs/etc/s6-overlay/s6-rc.d/cherry/run
#!/command/execlineb -P
with-contenv
s6-setuidgid nodejs
cd /app
node ./build/index.js
